# foodbot.service
[Unit]
Description=FoodBot Essensanmeldung
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%INSTALL_DIR%
Environment="PATH=/usr/bin:/usr/local/bin"
EnvironmentFile=%INSTALL_DIR%/.env
ExecStart=/usr/bin/python3 %INSTALL_DIR%/app/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
