# deployment/foodbot.service
[Unit]
Description=FoodBot - Feuerwehr Essensanmeldung
After=network.target

[Service]
Type=notify
User=%USER%
Group=%USER%
WorkingDirectory=%INSTALL_DIR%
Environment="PATH=%INSTALL_DIR%/venv/bin"
Environment="FLASK_ENV=production"
EnvironmentFile=%INSTALL_DIR%/.env
ExecStart=%INSTALL_DIR%/venv/bin/gunicorn -c %INSTALL_DIR%/deployment/gunicorn.conf.py "app:create_app()"
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateDevices=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=%INSTALL_DIR%/backups /var/log/foodbot

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=foodbot

[Install]
WantedBy=multi-user.target
